# services/tomebox.service
[Unit]
Description=TomeBox Audiobook Server
Documentation=https://github.com/yourusername/tomebox
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tomebox
Group=tomebox
WorkingDirectory=/opt/tomebox
ExecStart=/opt/tomebox/TomeBox --headless --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/tomebox

[Install]
WantedBy=multi-user.target